# systemd/bybit_xsreversal.service
[Unit]
Description=Bybit XS Reversal Bot (daily rebalance)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
#
# NOTE: This file is a TEMPLATE.
# The installer (`install.sh`) will generate a concrete unit at:
#   /etc/systemd/system/bybit_xsreversal.service
# with the correct WorkingDirectory/EnvironmentFile/ExecStart for your install path.
#
# If you install manually, set these to your actual PROJECT_DIR:
# WorkingDirectory=/home/ubuntu/cross-sectional-short-term-reversal/bybit_xsreversal
# EnvironmentFile=/home/ubuntu/cross-sectional-short-term-reversal/bybit_xsreversal/.env
# ExecStart=/home/ubuntu/cross-sectional-short-term-reversal/bybit_xsreversal/.venv/bin/bybit-xsreversal --config /home/ubuntu/cross-sectional-short-term-reversal/bybit_xsreversal/config/config.yaml live
WorkingDirectory=/opt/bybit_xsreversal/bybit_xsreversal
EnvironmentFile=/opt/bybit_xsreversal/bybit_xsreversal/.env
ExecStart=/opt/bybit_xsreversal/bybit_xsreversal/.venv/bin/bybit-xsreversal --config /opt/bybit_xsreversal/bybit_xsreversal/config/config.yaml live
Restart=on-failure
RestartSec=5

# Hardening (optional)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/bybit_xsreversal/bybit_xsreversal/outputs /opt/bybit_xsreversal/bybit_xsreversal/data_cache

[Install]
WantedBy=multi-user.target


